# Desktop/Python_cron/crontab
0 23 * * * root python /app/fetch_data.py
